# packaging/systemd/repovec-qdrant-api-key.service
[Unit]
Description=Provision the repovec Qdrant API key
Documentation=man:podman-secret(1)
Wants=systemd-sysusers.service
After=systemd-sysusers.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/repovec/repovec-qdrant-api-key
StandardOutput=journal
StandardError=journal
