# packaging/systemd/repovec-mcpd.service
[Unit]
Description=Repovec MCP bridge daemon
Requires=qdrant.service repovecd.service
After=qdrant.service repovecd.service

[Service]
User=repovec
Group=repovec
WorkingDirectory=/var/lib/repovec
Environment=HOME=/var/lib/repovec
ExecStart=/usr/bin/repovec-mcpd
